# config/node-exporter.service
[Unit]
Description=node-exporter

[Service]
ExecStart=/home/lzuccarelli/microservices/node_exporter
Restart=Always
PIDFile=/tmp/node-exporter_pid
